# crontab
0 0 1 10 * cd /; /usr/local/bin/node index.js
